# ext/ips/puppetmaster.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="puppet:master">

  <service name="network/puppet/master" type="service" version="1">
    <create_default_instance enabled="false"/>
    <single_instance/>

    <dependency name="config-file" grouping="require_all" restart_on="none" type="path">
      <service_fmri value="file:///etc/puppet/puppet.conf"/>
    </dependency>

    <dependency name='filesystem-local' grouping='require_all' restart_on='none' type='service'>
        <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>

    <dependency name='network' grouping='require_all' restart_on='error' type='service'>
        <service_fmri value='svc:/milestone/network:default'/>
    </dependency>

    <dependency name="fs-local" grouping="require_all" restart_on="none" type="service">
      <service_fmri value="svc:/system/filesystem/local"/>
    </dependency>

    <exec_method type="method" name="start" exec="/lib/svc/method/puppet-master start" timeout_seconds="60"/>

    <exec_method type="method" name="stop" exec="/lib/svc/method/puppet-master stop" timeout_seconds="60"/>

  <stability value="Evolving"/>

    <template>
      <common_name>
        <loctext xml:lang="C">Puppet Master Daemon</loctext>
      </common_name>
      <documentation>
        <manpage title="puppet" section="1"/>
        <doc_link name="puppetlabs.com" uri="https://puppetlabs.com/puppet/introduction"/>
      </documentation>
    </template>
  </service>

</service_bundle>
