# spec/fixtures/integration/provider/cron/crontab/crontab_user2
# HEADER: some simple
# HEADER: header
# Puppet Name: some_unrelevant job
* * * * * /bin/true
