# spec/fixtures/integration/provider/cron/crontab/moved_cronjob_input2
# HEADER: some simple
# HEADER: header
# Puppet Name: My daily failure
MAILTO=""
@daily /bin/false
# Puppet Name: some_unrelevant job
* * * * * /bin/true
