# spec/fixtures/integration/provider/cron/crontab/moved_cronjob_input2
# HEADER: some simple
# HEADER: header
# Puppet Name: some_unrelevant job
* * * * * /bin/true
# Puppet Name: My daily failure
@daily /bin/false
