# spec/fixtures/integration/provider/cron/crontab/purged
# HEADER: some simple
# HEADER: header

# commend with blankline above and below


# Puppet Name: only managed entry
* * * * * /bin/true
