# ext/solaris/smf/puppetmasterd.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<!-- PuppetMasterd Manifest: Luke Kanies - reductivelabs.com -->

<service_bundle type='manifest' name='puppetmasterd'>


<service
    name='network/puppetmasterd'
    type='service'
    version='1'>

    <create_default_instance enabled='true'/>
    <single_instance/>

    <dependency name='config-file'
                grouping='require_all'
                restart_on='none'
                type='path'>
        <service_fmri value='file:///etc/puppet/puppet.conf'/>
    </dependency>

    <dependency name='loopback'
                grouping='require_all'
                restart_on='error'
                type='service'>
        <service_fmri value='svc:/network/loopback:default'/>
    </dependency>

    <dependency name='physical'
                grouping='require_all'
                restart_on='error'
                type='service'>
        <service_fmri value='svc:/network/physical:default'/>
    </dependency>

    <dependency name='fs-local'
                grouping='require_all'
                restart_on='none'
                type='service'>
        <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>

    <exec_method
        type='method'
        name='start'
        exec='/opt/csw/lib/svc/method/svc-puppetmasterd start'
        timeout_seconds='60' />

    <exec_method
        type='method'
        name='stop'
        exec='/opt/csw/lib/svc/method/svc-puppetmasterd stop'
        timeout_seconds='60' />

    <exec_method
        type='method'
        name='refresh'
        exec='/opt/csw/lib/svc/method/svc-puppetmasterd reload'
        timeout_seconds='60' />

    <stability value='Unstable' />

    <template>
        <common_name>
            <loctext xml:lang='C'>Puppet Server Daemon</loctext>
        </common_name>
        <documentation>
            <manpage title='puppetmasterd' section='1' />
            <doc_link name='reductivelabs.com'
                uri='http://www.reductivelabs.com/projects/puppet' />
        </documentation>
    </template>
</service>

</service_bundle>
